# config/cron/drupal.crontab
0 * * * * sudo -u www-data sh -c "cd /opt/xtuple/commerce/live && ./console.php cron:run" >> /var/log/xtuple/cron/live.log 2>&1
15 * * * * sudo -u www-data sh -c "cd /opt/xtuple/commerce/stage && ./console.php cron:run" >> /var/log/xtuple/cron/stage.log 2>&1
30 * * * * sudo -u www-data sh -c "cd /opt/xtuple/commerce/dev && ./console.php cron:run" >> /var/log/xtuple/cron/dev.log 2>&1
